# crontab
SHELL=/bin/bash
# Run the script every 5 minutes
12,32,52 * * * *  python /app/bfx_rate.py >> /var/log/cron.log 2>&1